# supervisord.conf
[supervisord]
nodaemon=true

[program:php]
command=/usr/sbin/php-fpm7.0 --nodaemonize

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

